# init.d-caps2esclogo4alt
#!/sbin/openrc-run
# install to /etc/init.d, ln -s to /etc/runtimes/default (rc-update enable)
name="caps2esc-logo4alt"
supervisor="supervise-daemon"
command="-N -20 /usr/local/bin/udevmon"
pidfile="/run/${RC_SVCNAME}.pid"
command_args="-c /etc/caps2esclogo4alt.yaml"

reload() {
    ebegin "Reloading caps2esc-logo4alt configuration"
    ${supervisor} ${RC_SVCNAME} --signal HUP --pidfile "${pidfile}"
    eend $?
}
